# packaging/linux-common/automator-scheduler.timer
[Unit]
Description=Run the Automator scheduler tick every minute

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1sec
Persistent=true

[Install]
WantedBy=timers.target
